# ffff.conf
[program:ffff]
command=python3 /home/curator2/ffff/main.py
autostart=true
autorestart=true
stderr_logfile=/home/curator2/ffff/test.err.log
stdout_logfile=/home/curator2/ffff/test.out.log